# deployment/systemd/backend.service
[Unit]
Description=totem backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=DEPLOY_USER
WorkingDirectory=PROJECT_ROOT/apps/backend
Environment="NODE_ENV=production"
Environment="PATH=DEPLOY_HOME/.bun/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=DEPLOY_HOME/.bun/bin/bun run --env-file=PROJECT_ROOT/.env.production src/index.ts

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=PROJECT_ROOT/data

MemoryMax=2G
TasksMax=512

StandardOutput=journal
StandardError=journal
SyslogIdentifier=totem-backend

[Install]
WantedBy=multi-user.target
