# deployment/systemd/frontend.service
[Unit]
Description=totem frontend
After=network-online.target totem-backend.service
Wants=network-online.target totem-backend.service

[Service]
Type=simple
User=DEPLOY_USER
WorkingDirectory=PROJECT_ROOT/apps/frontend
Environment="NODE_ENV=production"
Environment="PATH=DEPLOY_HOME/.bun/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=DEPLOY_HOME/.bun/bin/bun --bun vite preview --host 0.0.0.0 --port 5173

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

MemoryMax=1G
TasksMax=256

StandardOutput=journal
StandardError=journal
SyslogIdentifier=totem-frontend

[Install]
WantedBy=multi-user.target
